# install/mulch-proxy.sample.service
[Unit]
Description=Mulch Reverse Proxy
After=network-online.target

[Service]
User={USER}
ExecStart={MULCH_PROXY}
Type=simple
Restart=on-failure
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
